# systemd/gitSync.service
[Unit]
Description="Sync with GitHub"


[Service]
ExecStart=/bin/bash ${HOME}/Documents/scripts/allScripts/gitSync.sh
